# configs/etc/profile.d/default-editor.sh
## set default editor

if [ -z "$EDITOR" ]; then
  export EDITOR="/usr/bin/nano"
fi